# .crontab
0 * * * * php bin/console app:application:command:update-status